# 10-usbdriveadded.rules
KERNEL=="sd?1", ACTION=="add", RUN+="/home/pi/startstartloop.sh" 
